# debian/courier-imap.service
[Unit]
Description=Courier IMAP Daemon
Requires=courier-authdaemon.service
After=courier-authdaemon.service

[Service]
Type=forking
ExecStart=/usr/sbin/imapd start
ExecStop=/usr/sbin/imapd stop
PIDFile=/run/courier/imapd.pid

[Install]
WantedBy=multi-user.target
